# goagent.conf
description "start and stop the goagent"
version "2.0"
author "bob zhang"

start on runlevel [2345]
stop on runlevel [!2345]
script
  /usr/bin/python /home/software/goagent/local/proxy.py
end script
